# simple.conf
[program:myserver]
command=gunicorn -w 4 -b 0.0.0.0:80 PicSift:app

[supervisord]
logfile=/home/ubuntu/PicSiftApp/supervisord.log
loglevel=debug
user=root
